# .config/autostart/sxhkd.desktop
[Desktop Entry]
Name=sxhkd
Comment=Simple X hotkey daemon
Exec=/run/current-system/sw/bin/sxhkd
Terminal=false
Type=Application
